# systemd/edgewatch-ble.service
[Unit]
Description=EdgeWatch Ble Service
After=network.targetedgewatch-zmq-proxy.service
Requires=edgewatch-zmq-proxy.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/edgewatch
ExecStart=/usr/bin/python3 src/services/ble_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
